# ubuntu/poolmon.conf
description "Poolmon Dovecot Director Monitoring Daemon"

start on startup
stop on runlevel [016]

kill timeout 3
expect fork

respawn

exec poolmon -l syslog -c /etc/poolmon-credentials.conf
